# software/systemd/lego-sorter-ui.service
[Unit]
Description=LEGO Sorter UI
Documentation=https://basicallysource.github.io/sorter-v2/sorter/installation/
After=network-online.target lego-sorter-backend.service
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__SOFTWARE_DIR__/ui
ExecStart=__PNPM_BIN__ preview --host 0.0.0.0 --port 5173
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
