# software/systemd/sorter-ui-dev.service
[Unit]
Description=Sorter UI (dev — Vite HMR)
After=network-online.target sorter-backend-dev.service
Wants=network-online.target
Conflicts=sorter-ui.service
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=simple
User=__USER__
WorkingDirectory=__SOFTWARE_DIR__/sorter/frontend
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# `vite dev` gives sub-second HMR for .svelte / .ts edits scp'd from the Mac.
# `--host 0.0.0.0` makes it reachable on the LAN; `--strictPort` makes it fail
# loudly if the port is busy instead of silently picking another port.
# Port 80 lets users hit http://<host>.local without a port suffix; unprivileged
# binding is enabled via /etc/sysctl.d/99-sorter-unprivileged-ports.conf.
ExecStart=__PNPM_BIN__ dev --host 0.0.0.0 --port 80 --strictPort
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
