# software/sorteros/build/overlay/etc/profile.d/99-sorter-aliases.sh
# Default shell aliases baked into SorterOS. Sourced by login shells (e.g. SSH
# sessions) for every user, including the pre-existing orangepi account that
# /etc/skel never touches.
alias ll='ls -la'
