# software/sorteros/build/overlay/etc/systemd/system/sorteros-bootstrap-users.service
[Unit]
Description=SorterOS bootstrap users and SSH homes
DefaultDependencies=no
After=local-fs.target
Before=ssh.service sorteros-firstboot.service multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/sorteros-bootstrap-users.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
