# software/sorteros/build/overlay/etc/systemd/system/sorter-backend-dev.service
[Unit]
Description=Sorter Backend (dev)
After=network-online.target
Wants=network-online.target
Conflicts=sorter-backend.service
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=simple
User=__USER__
WorkingDirectory=__SOFTWARE_DIR__/sorter/backend
# Dev mode notes:
#   - Same supervisor + main.py as prod so hardware/coordinator behavior matches.
#   - For fast iteration on a code edit, prefer a soft restart of main.py via the
#     supervisor control endpoint instead of `systemctl restart`:
#       curl -sS -X POST http://127.0.0.1:8001/api/supervisor/restart
#     That skips systemd's RestartSec=5 and the supervisor warmup.
ExecStart=/bin/bash -c "set -a; source __SOFTWARE_DIR__/.env; exec __UV_BIN__ run python supervisor.py"
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
